# web-geo-server/supervisord.ini
[program:opendai-django]
directory=/var/www/odai_pilots
user=root
command=gunicorn_script.sh
stdout_logfile=var/log/opendai/django_stdout.log
stdout_logfile_maxbytes=100MB

[program:openday-cerely]
directory=/var/www/odai_pilots
user=root
command=celery_script.sh
stdout_logfile=/var/log/opendai/celery_stdout.log
stdout_logfile_maxbytes=100MB